# debian/vpn-user-portal.cron.d
@hourly www-data /usr/bin/vpn-user-portal-foreign-key-list-fetcher
